# glances.service
[Unit]
Description=Glances

[Service]
ExecStart=/scripts/glances-websrv
Restart=on-abort

[Install]
WantedBy=multi-user.target
